# mcast-pa/files/etc/uci-defaults/mcast-pa
uci -q get network.lan && {
	uci -q set network.lan.igmp_snooping=1
	uci commit
}
