# lq-wave-300_5_x/files/etc/hotplug.d/net/00-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && [ "${INTERFACE%%[0-9]}" = "wlan" ] && {
	/sbin/wifi config
}
